# stage-piproxy/01-configure/files/oled-status.service
[Unit]
Description=OLED network status display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/oled-status.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
